# warp-room.service
[Unit]
Description=Warp-as-Room Tile Classifier — CPU simulation of GPU room inference
Documentation=https://github.com/Lucineer/edge-llama

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/lucineer/edge-llama/warp_room.py --train
Environment=PLATO_URL=http://localhost:8847
User=lucineer
